# bootstrap/var/jb/Library/dpkg/info/shshd.list
/.
/var
/var/jb
/var/jb/Library
/var/jb/Library/LaunchDaemons
/var/jb/Library/LaunchDaemons/us.diatr.shshd.plist
/var/jb/usr
/var/jb/usr/libexec
/var/jb/usr/libexec/shshd-wrapper
/var/jb/usr/sbin
/var/jb/usr/sbin/shshd
/var/jb/usr/share
/var/jb/usr/share/doc
/var/jb/usr/share/doc/shshd
/var/jb/usr/share/doc/shshd/LICENSE
