# bootstrap/var/jb/Library/dpkg/info/system-cmds.list
/.
/var
/var/jb
/var/jb/Library
/var/jb/Library/LaunchDaemons
/var/jb/Library/LaunchDaemons/com.apple.atrun.plist
/var/jb/bin
/var/jb/bin/sync
/var/jb/bin/wait4path
/var/jb/etc
/var/jb/etc/pam.d
/var/jb/etc/pam.d/login
/var/jb/etc/pam.d/login.term
/var/jb/etc/pam.d/passwd
/var/jb/sbin
/var/jb/sbin/dmesg
/var/jb/sbin/dynamic_pager
/var/jb/sbin/nologin
/var/jb/sbin/reboot
/var/jb/sbin/shutdown
/var/jb/usr
/var/jb/usr/bin
/var/jb/usr/bin/arch
/var/jb/usr/bin/at
/var/jb/usr/bin/chpass
/var/jb/usr/bin/cpuctl
/var/jb/usr/bin/fs_usage
/var/jb/usr/bin/getconf
/var/jb/usr/bin/hostinfo
/var/jb/usr/bin/latency
/var/jb/usr/bin/login
/var/jb/usr/bin/lskq
/var/jb/usr/bin/lsmp
/var/jb/usr/bin/ltop
/var/jb/usr/bin/memory_pressure
/var/jb/usr/bin/newgrp
/var/jb/usr/bin/pagesize
/var/jb/usr/bin/passwd
/var/jb/usr/bin/purge
/var/jb/usr/bin/sc_usage
/var/jb/usr/bin/stackshot
/var/jb/usr/bin/trace
/var/jb/usr/libexec
/var/jb/usr/libexec/atrun
/var/jb/usr/libexec/getty
/var/jb/usr/sbin
/var/jb/usr/sbin/ac
/var/jb/usr/sbin/accton
/var/jb/usr/sbin/adduser
/var/jb/usr/sbin/iostat
/var/jb/usr/sbin/mkfile
/var/jb/usr/sbin/pw
/var/jb/usr/sbin/pwd_mkdb
/var/jb/usr/sbin/rmuser
/var/jb/usr/sbin/sysctl
/var/jb/usr/sbin/taskpolicy
/var/jb/usr/sbin/vifs
/var/jb/usr/sbin/vipw
/var/jb/usr/sbin/zdump
/var/jb/usr/sbin/zic
/var/jb/usr/share
/var/jb/usr/share/man
/var/jb/usr/share/man/man1
/var/jb/usr/share/man/man1/arch.1.zst
/var/jb/usr/share/man/man1/at.1.zst
/var/jb/usr/share/man/man1/chpass.1.zst
/var/jb/usr/share/man/man1/fs_usage.1.zst
/var/jb/usr/share/man/man1/getconf.1.zst
/var/jb/usr/share/man/man1/latency.1.zst
/var/jb/usr/share/man/man1/login.1.zst
/var/jb/usr/share/man/man1/lskq.1.zst
/var/jb/usr/share/man/man1/lsmp.1.zst
/var/jb/usr/share/man/man1/ltop.1.zst
/var/jb/usr/share/man/man1/machine.1.zst
/var/jb/usr/share/man/man1/memory_pressure.1.zst
/var/jb/usr/share/man/man1/newgrp.1.zst
/var/jb/usr/share/man/man1/pagesize.1.zst
/var/jb/usr/share/man/man1/passwd.1.zst
/var/jb/usr/share/man/man1/trace.1.zst
/var/jb/usr/share/man/man1/vm_stat.1.zst
/var/jb/usr/share/man/man1/wait4path.1.zst
/var/jb/usr/share/man/man1/zprint.1.zst
/var/jb/usr/share/man/man5
/var/jb/usr/share/man/man5/adduser.conf.5.zst
/var/jb/usr/share/man/man5/gettytab.5.zst
/var/jb/usr/share/man/man5/nologin.5.zst
/var/jb/usr/share/man/man5/pw.conf.5.zst
/var/jb/usr/share/man/man5/sysctl.conf.5.zst
/var/jb/usr/share/man/man5/ttys.5.zst
/var/jb/usr/share/man/man8
/var/jb/usr/share/man/man8/ac.8.zst
/var/jb/usr/share/man/man8/accton.8.zst
/var/jb/usr/share/man/man8/adduser.8.zst
/var/jb/usr/share/man/man8/atrun.8.zst
/var/jb/usr/share/man/man8/cpuctl.8.zst
/var/jb/usr/share/man/man8/dmesg.8.zst
/var/jb/usr/share/man/man8/dynamic_pager.8.zst
/var/jb/usr/share/man/man8/getty.8.zst
/var/jb/usr/share/man/man8/hostinfo.8.zst
/var/jb/usr/share/man/man8/iostat.8.zst
/var/jb/usr/share/man/man8/mkfile.8.zst
/var/jb/usr/share/man/man8/nologin.8.zst
/var/jb/usr/share/man/man8/nvram.8.zst
/var/jb/usr/share/man/man8/purge.8.zst
/var/jb/usr/share/man/man8/pw.8.zst
/var/jb/usr/share/man/man8/pwd_mkdb.8.zst
/var/jb/usr/share/man/man8/reboot.8.zst
/var/jb/usr/share/man/man8/rmuser.8.zst
/var/jb/usr/share/man/man8/sa.8.zst
/var/jb/usr/share/man/man8/shutdown.8.zst
/var/jb/usr/share/man/man8/sync.8.zst
/var/jb/usr/share/man/man8/sysctl.8.zst
/var/jb/usr/share/man/man8/taskpolicy.8.zst
/var/jb/usr/share/man/man8/vifs.8.zst
/var/jb/usr/share/man/man8/vipw.8.zst
/var/jb/usr/share/man/man8/zdump.8.zst
/var/jb/usr/share/man/man8/zic.8.zst
/var/jb/usr/share/skel
/var/jb/sbin/halt -> /var/jb/sbin/reboot
/var/jb/usr/bin/chfn -> chpass
/var/jb/usr/bin/chsh -> chpass
/var/jb/usr/bin/machine -> arch
/var/jb/usr/sbin/nologin -> /var/jb/sbin/nologin
/var/jb/usr/sbin/reboot -> /var/jb/sbin/reboot
/var/jb/usr/share/man/man1/chfn.1.zst -> chpass.1.zst
/var/jb/usr/share/man/man1/chsh.1.zst -> chpass.1.zst
/var/jb/usr/share/man/man8/halt.8.zst -> reboot.8.zst
